# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
directory=/app
command=python doge.py

[program:server]
directory=/app
command=python server.py

[supervisorctl]